# src/update/modules/500-jellyfin/root/usr/lib/systemd/system/jellyfin.service
[Unit]
After=network-online.target

[Service]
ExecStart=/usr/bin/jellyfin

[Install]
RequiredBy=multi-user.target
