# src/update/modules/500-jellyfin/root/usr/lib/systemd/system/jellyfin.service.d/dirs.conf
[Service]
Environment="JELLYFIN_DATA_DIR=/mnt/fossil/config/jellyfin/data"
Environment="JELLYFIN_CONFIG_DIR=/mnt/fossil/config/jellyfin/config"
Environment="JELLYFIN_CACHE_DIR=/var/cache/jellyfin"
Environment="JELLYFIN_LOG_DIR=/var/log/jellyfin"
