# src/modules/500-build-fstab/root/usr/lib/systemd/system/mb-gen-fstab.service
[Unit]
After=mnt-fossil.mount
Requires=mnt-fossil.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/mb-gen-fstab

[Install]
WantedBy=multi-user.target
