# src/modules/400-x11vnc/root/usr/lib/systemd/system/x11vnc.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/x11vnc \
        -display :0 \
        -no6 \
        -rfbport 5901 \
        -nevershared \
        -forever \
        -auth /var/run/lightdm/root/:0

