# src/modules/500-steam/root/usr/lib/systemd/system/steam-link-download-dir.service
[Service]
Type=oneshot
ExecStart=/usr/bin/setup-steam-library-symlink

[Install]
WantedBy=graphical.target
